# etc/supervisor/groups/monadical.homenet.conf
[group:monadical.homenet]
programs=
    penny-postgres,
    penny-django,
    penny-nginx

[program:penny-argo]
autostart=true
